# deploy/sablos.service
[Unit]
Description=sablos Django site (gunicorn)
After=network.target

[Service]
Type=notify
User=sablos
Group=sablos
WorkingDirectory=/var/www/sablos_root
EnvironmentFile=/var/www/sablos_root/.env
ExecStart=/var/www/sablos_root/.venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/run/sablos/sablos.sock \
    --access-logfile - \
    --error-logfile - \
    config.wsgi:application
Restart=on-failure
RestartSec=5
RuntimeDirectory=sablos
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
